# ASSETO/asseto/asseto.service
[Unit]
Description=ASSETO Inventory Management
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/asseto
EnvironmentFile=/opt/asseto/.env
ExecStart=/opt/asseto/venv/bin/gunicorn \
    --bind 127.0.0.1:5000 \
    --workers 3 \
    --threads 2 \
    --timeout 60 \
    --access-logfile /opt/asseto/logs/access.log \
    --error-logfile /opt/asseto/logs/error.log \
    app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5s
KillMode=mixed
TimeoutStopSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/asseto/logs /opt/asseto/backups /opt/asseto/static/photos /opt/asseto/static/signatures

[Install]
WantedBy=multi-user.target
